# scripts/configs/supervisor.ini
[program:yurika]
directory = /opt/yurika
command = /opt/yurika/.venv/bin/gunicorn --env DJANGO_SETTINGS_MODULE=project.settings --pythonpath /opt/yurika/ --bind=unix:/opt/yurika/gunicorn.sock project.wsgi
stdout_logfile = /opt/yurika/stdout.log
redirect_stderr = true
autostart = true
autorestart = true
user = mortar
group = mortar

[program:celery]
directory = /opt/yurika
command = /opt/yurika/.venv/bin/celery -A project.celery worker -l info
stdout_logfile = /opt/yurika/celery.log
redirect_stderr = true
autostart = true
autorestart = true
user = mortar
group = mortar
